# templates/etc/profile.d/98-hashicorp-packer.j2
PACKER_CACHE_DIR={{ hashicorp_packer_cache_dir }}
PACKER_LOG={{ hashicorp_packer_log }}
PACKER_LOG_PATH={{ hashicorp_packer_log_path }}
CHECKPOINT_DISABLE={{ hashicorp_packer_checkpoint_disable }}
